# playbooks/files/90-mikes-udev.rules
ACTION=="add", SUBSYSTEM=="bluetooth", ATTRS{address}=="bc:d0:74:b7:f7:5f", RUN+="/usr/bin/bluetoothctl connect bc:d0:74:b7:f7:5f"
